# rpi_example/tps.service
[Unit]
Description="syncs podcasts to tonies"

[Service]
ExecStart=/opt/tps/tps.sh